# systemd/rmkit-cn-upload.service
[Unit]
Description=rmkit-cn 文件上传服务
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/root/rmkit-cn/upload-server
ExecStart=/home/root/rmkit-cn/upload-server/upload-server \
  -listen 0.0.0.0:8080 \
  -static /home/root/rmkit-cn/upload-server/static \
  -fonts /home/root/.local/share/rmkit-cn/fonts \
  -screens /home/root/.local/share/rmkit-cn/screens \
  -staging /tmp/rmkit_upload \
  -fonts-active /home/root/.local/share/fonts \
  -xochitl-conf /home/root/.config/remarkable/xochitl.conf
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
